# runner/crontab
0 4 * * * pennmush /home/pennmush/backup/pennback.sh
