# runner/crontab
0 4 * * * pennmush sh /home/pennmush/backup/penback.sh
